# docker-binaries/docker-manager.timer
[Unit]
Description=Cleanup Docker Daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=docker-manager.service

[Install]
WantedBy=timers.target
